# deployment/systemd/bostad-backend.service
[Unit]
Description=Bostad FastAPI backend
After=network-online.target
Wants=network-online.target

[Service]
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_DIR__/backend
EnvironmentFile=-__APP_DIR__/deployment/config.env
Environment=PYTHONUNBUFFERED=1
Environment=LOG_LEVEL=INFO
ExecStart=__APP_DIR__/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
