# deployment/systemd/bostad-refetch.service
[Unit]
Description=Fetch and cache all bostad listings
After=network-online.target bostad-backend.service
Wants=network-online.target

[Service]
Type=oneshot
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_DIR__
EnvironmentFile=-__APP_DIR__/deployment/config.env
ExecStart=__BIN_DIR__/refetch-listings.sh
